# deploy/supervisor/smart-lms-ai-scan.conf
[program:smart-lms-ai-scan]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/smart_lms/artisan queue:work database --queue=ai-scan --tries=1 --sleep=1 --timeout=600 --backoff=5 --memory=768
directory=/var/www/smart_lms
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/smart_lms/storage/logs/ai-scan-worker.log
stopwaitsecs=610

